# spectrum/image/rootfs/usr/lib/systemd/system/dbus-broker.service.d/60-spectrum-desktop-protection.conf
[Service]
MemoryMin=64M
MemoryLow=128M
CPUWeight=1000
IOWeight=1000
ManagedOOMPreference=avoid
