# spectrum/image/rootfs/usr/lib/systemd/user/org.gnome.Shell@.service.d/60-spectrum-desktop-protection.conf
[Service]
CPUWeight=1000
IOWeight=1000
ManagedOOMPreference=avoid
