# spectrum/image/rootfs/usr/lib/systemd/user/org.gnome.SettingsDaemon.XSettings.service.d/60-spectrum-desktop-protection.conf
[Service]
MemoryMin=32M
MemoryLow=64M
CPUWeight=800
IOWeight=800
ManagedOOMPreference=avoid
