# spectrum/image/rootfs/usr/lib/systemd/user/xdg-permission-store.service.d/60-spectrum-desktop-protection.conf
[Service]
MemoryMin=64M
MemoryLow=128M
CPUWeight=800
IOWeight=800
ManagedOOMPreference=avoid
