# spectrum/image/rootfs/usr/lib/systemd/system/kanata-main.service.d/60-spectrum-desktop-protection.conf
[Service]
MemoryMin=32M
MemoryLow=64M
CPUWeight=1000
IOWeight=1000
ManagedOOMPreference=avoid
